# deploy/kepagent.service
[Unit]
Description=KepAgent node control agent
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/kepagent
ExecStart=/opt/kepagent/kepagent --config /opt/kepagent/config.yaml
Restart=always
RestartSec=5
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=multi-user.target
